# dotfiles/.config/systemd/user/headset-mic-monitor.service
[Unit]
Description=Auto-switch microphone on headphone jack events
After=pipewire.service wireplumber.service

[Service]
ExecStart=%h/.local/bin/headset-mic-monitor
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
